# deploy/ai-os-hfm.service
[Unit]
Description=AI OS — HFM (Holistic Functional Care) on :8788
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-os
EnvironmentFile=/opt/ai-os/.env.hfm
Environment=NODE_ENV=production
Environment=SERVE_STATIC=true
ExecStart=/usr/bin/node /opt/ai-os/node_modules/.bin/tsx server/index.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-os-hfm

[Install]
WantedBy=multi-user.target
