# deploy/ai-os-sc.service
[Unit]
Description=AI OS — SC (Simple Connect) on :8787
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-os
EnvironmentFile=/opt/ai-os/.env.sc
Environment=NODE_ENV=production
Environment=SERVE_STATIC=true
ExecStart=/usr/bin/node /opt/ai-os/node_modules/.bin/tsx server/index.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-os-sc

[Install]
WantedBy=multi-user.target
